# src/forecast/crons/realtime.prd.cron
#crontab on cpcwork2 is UTC time
MAILTO=peitao.peng@noaa.gov

# CA week3-4 forecst
58 12 * * * /cpc/home/wd52pp/project/ca_prd/wk34/opr/realtime.run.ca_wk34.sh > /cpc/home/wd52pp/data/wk34_fcst/ca/0log1 2>&1
#25 13 * * * /cpc/home/wd52pp/project/ca_prd/wk34/opr/realtime.run.ca_wk34.mwic.sh > /cpc/home/wd52pp/data/wk34_fcst/ca/0log1 2>&1
# submonthly and monthly forecast
00 16 5 * * bash /home/ppeng/data/downloads/fetch_ftp.bash > /home/ppeng/data/downloads/0log1 2>&1
00 17 5 * * bash /home/ppeng/src/develop/data_proc/fetch_ftp.bash > /home/ppeng/data/downloads/0log1 2>&1
55 15 * * * /cpc/home/wd52pp/project/consistency/opr/realtime.run.gefs_data.sh > /cpc/consistency/gefs_new/0log1 2>&1
05 16 * * * /cpc/home/wd52pp/project/consistency/opr/realtime.run.GEFS.subssfcst.sh > /cpc/consistency/subss_fcst/0log1 2>&1
25 16 * * * /cpc/home/wd52pp/project/consistency/opr/realtime.run.copy_2_web.sh > /cpc/consistency/subss_fcst/0log1 2>&1
# CA seasonal forecst
25 8 5 * * /cpc/home/wd52pp/project/ca_prd/sst/opr/realtime.run.ca_season.sst.sh> /cpc/home/wd52pp/data/season_fcst/ca/0log1 2>&1 
25 13 7 * * /cpc/home/wd52pp/project/ca_prd/sst/opr/realtime.run.ca_season.tpz.sh> /cpc/home/wd52pp/data/season_fcst/ca/0log1 2>&1 
# CA monthly forecst
25 8 6 * * /cpc/home/wd52pp/project/ca_prd/sst/opr/realtime.run.ca_month.sst.sh> /cpc/home/wd52pp/data/season_fcst/ca/0log1 2>&1 
25 8 8 * * /cpc/home/wd52pp/project/ca_prd/sst/opr/realtime.run.ca_month.tpz.sh> /cpc/home/wd52pp/data/season_fcst/ca/0log1 2>&1 

25 8 9 * * /cpc/home/wd52pp/project/ca_prd/sst/opr/final.run.mon_sst.sh> /cpc/home/wd52pp/data/season_fcst/ca/0log1 2>&1 
00 14 25 * * /cpc/home/wd52pp/project/ca_prd/sst/opr/prelim.run.mon_sst.sh> /cpc/home/wd52pp/data/season_fcst/ca/0log1 2>&1 
# verification for NMME past month and season
00 7 10 * * /cpc/home/wd52pp/project/id/vfc/src/realtime.nmme.HSS.scr> /cpc/consistency/id/nmme/skill/0log1 2>&1
